# detections/sigma/fake_service_remote_endpoint.yml
title: Fake Windows Service with Remote Endpoint Parameter
id: a1b2c3d4-e5f6-7890-abcd-ef1234567890
status: stable
description: >
  Detects a Windows service whose binary path contains a remote IP address or
  hostname as a command-line parameter. Legitimate Windows services do not
  accept remote endpoint parameters at service-start time. Observed in
  CRIMSON OSPREY: msadvapi2_64.exe and msadvapi2_32.exe registered as services
  with the C2 IP 172.16.6.12 as a hard-coded argument, enabling persistent
  F-Response remote access on every system boot.
references:
  - https://attack.mitre.org/techniques/T1543/003/
  - https://www.f-response.com/
author: CaseFile — SANS Find Evil! 2026
date: 2026-05-19
modified: 2026-05-19
tags:
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1543.003
logsource:
  product: windows
  category: registry_add
  definition: >
    Requires Sysmon EventID 13 (RegistryValueSet) or equivalent registry
    monitoring on HKLM\SYSTEM\CurrentControlSet\Services.
detection:
  selection:
    EventID: 13
    TargetObject|contains: '\SYSTEM\CurrentControlSet\Services\'
    TargetObject|endswith: '\ImagePath'
    Details|re: '(\d{1,3}\.){3}\d{1,3}'
  filter_legitimate:
    Details|contains:
      - 'svchost.exe'
      - 'lsass.exe'
      - 'services.exe'
  condition: selection and not filter_legitimate
falsepositives:
  - Remote monitoring agents that legitimately embed endpoint IPs in service paths
  - Third-party backup agents with site-specific configuration
level: high
fields:
  - TargetObject
  - Details
  - Image
  - ProcessId
